# deployment/supervisor/frameapp_supervisor.ini
[program:frameapp]
command=pipenv run gunicorn -b localhost:8000 -w 4 frameapp:app
directory=/home/fedora/projects/Flask_FrameApp
user=fedora
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
